# scripts/backup/crontab
# Backup container crontab (dcron format, root user).
#
# Installed in the image at /etc/crontabs/root (Dockerfile.backup).
# dcron runs as `crond -f` (default log level 5/NOTICE) so scheduled
# wakeups appear in `docker compose logs backup`.
#
# Times are Europe/Berlin — the container sets TZ=Europe/Berlin (see
# docker-compose.yml services.backup.environment and the tzdata package
# in Dockerfile.backup). DST transitions are handled by tzdata; no
# manual adjustment needed twice a year.
#
# Schedule (Europe/Berlin local time):
#   09:00, 12:00, 15:00, 18:00, 21:00  weekdays  run-backup
#   12:00                              weekends  run-backup
#   +2 min offset on each tick                   run-drill
#
# Rationale: humans edit data during working hours; a 3-hour backup
# window matches the pace at which business data actually changes.
# Drill follows the same cadence so every backup is Tier-2 verified
# when the operator's key is loaded (AC-168 skip path when absent).
#
# No rotation entry: retention is linear (14-day bucket lock + 90-day
# lifecycle) — R2 handles the deletes. See ADR-0020 §Decision and
# docs/ops/backup/overview.md.
#
# Output redirection (`> /proc/1/fd/1 2>/proc/1/fd/2`): dcron's default
# behaviour is to mail a job's stdout/stderr via sendmail. This image
# carries no MTA, so without redirection the script output is silently
# dropped — leaving the operator with only the dcron NOTICE wakeup line
# and no per-run "ok" / "skipped" / "failed" cue from the runner.
# /proc/1/fd/{1,2} are crond's own (PID 1) stdout/stderr, which Docker's
# json-file driver captures.

0 9,12,15,18,21 * * 1-5 /usr/local/bin/run-backup.sh > /proc/1/fd/1 2>/proc/1/fd/2
0 12            * * 6,0 /usr/local/bin/run-backup.sh > /proc/1/fd/1 2>/proc/1/fd/2
2 9,12,15,18,21 * * 1-5 /usr/local/bin/run-drill.sh  > /proc/1/fd/1 2>/proc/1/fd/2
2 12            * * 6,0 /usr/local/bin/run-drill.sh  > /proc/1/fd/1 2>/proc/1/fd/2
